# deploy/bf_planning.service
[Unit]
Description=BF Planning duty schedule calendar
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/bf_planning
EnvironmentFile=/opt/bf_planning/.env
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
